# examples/systemd/autodns-backup-verify.timer
[Unit]
Description=Monthly AutoDNS backup verification

[Timer]
OnCalendar=monthly
Persistent=true

[Install]
WantedBy=timers.target
